# printa-client/etc/xdg/autostart/printa-client.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Printa Client
Comment=Printa Client GUI
Exec=/usr/bin/printa-client
Terminal=false
Type=Application
